# eudaemon.service
[Unit]
Description=Eudaemon
Wants=default.target

[Service]
Type=simple
ExecStart=eudaemon

[Install]
WantedBy=default.target
